# MojangPriangan/deploy/asm-httpd.service
# deploy/asm-httpd.service
[Unit]
Description=asmhttpd (NASM) demo
After=network.target

[Service]
ExecStart=/opt/asmhttpd/asmhttpd -p 8080 -d /opt/asmhttpd/www
User=www-data
Restart=always
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
